# deploy/systemd/rexfinhub-bloomberg.timer
[Unit]
Description=REX FinHub Bloomberg Pull Timer (every 15 min, 4-6 PM ET weekdays)

[Timer]
# VPS timezone set to America/New_York — DST handled automatically
OnCalendar=Mon..Fri 16:00:00
OnCalendar=Mon..Fri 16:15:00
OnCalendar=Mon..Fri 16:30:00
OnCalendar=Mon..Fri 16:45:00
OnCalendar=Mon..Fri 17:00:00
OnCalendar=Mon..Fri 17:15:00
OnCalendar=Mon..Fri 17:30:00
OnCalendar=Mon..Fri 17:45:00
Persistent=true

[Install]
WantedBy=timers.target
