# deploy/systemd/rexfinhub-reconciler.service
[Unit]
Description=rexfinhub daily SEC index reconciler
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=jarvis
WorkingDirectory=/home/jarvis/rexfinhub
EnvironmentFile=/home/jarvis/rexfinhub/config/.env
Environment=SEC_CACHE_DIR=/home/jarvis/rexfinhub/cache/sec
ExecStart=/home/jarvis/venv/bin/python -m etp_tracker.reconciler
TimeoutStartSec=600
StandardOutput=journal
StandardError=journal
