# deploy/systemd/rexfinhub-sec-scrape.service
[Unit]
Description=REX FinHub SEC Filing Scrape
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=jarvis
WorkingDirectory=/home/jarvis/rexfinhub
EnvironmentFile=/home/jarvis/rexfinhub/config/.env
# Pin process timezone to ET (audit fix R7).
Environment=TZ=America/New_York
Environment=SEC_CACHE_DIR=/home/jarvis/rexfinhub/cache/sec
ExecStart=/home/jarvis/venv/bin/python scripts/run_all_pipelines.py --skip-email --skip-market
TimeoutStartSec=7200
StandardOutput=journal
StandardError=journal
