# deploy/systemd/rexfinhub-daily.service
[Unit]
Description=REX FinHub Daily Pipeline (full sync + email reports)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=jarvis
WorkingDirectory=/home/jarvis/rexfinhub
EnvironmentFile=/home/jarvis/rexfinhub/config/.env
Environment=SEC_CACHE_DIR=/home/jarvis/rexfinhub/cache/sec
ExecStart=/home/jarvis/venv/bin/python scripts/run_daily.py --skip-sec
TimeoutStartSec=3600
StandardOutput=journal
StandardError=journal
