# deploy/systemd/rexfinhub-daily.timer
[Unit]
Description=REX FinHub Daily Pipeline Timer (weekdays 5:00 PM ET)

[Timer]
# VPS timezone set to America/New_York — DST handled automatically
OnCalendar=Mon..Fri 17:00:00
Persistent=true

[Install]
WantedBy=timers.target
