# deploy/systemd/rexfinhub-daily.timer
[Unit]
Description=REX FinHub Daily Pipeline Timer (weekdays 19:30 ET)

[Timer]
# VPS timezone set to America/New_York — DST handled automatically.
# 2026-04-28 cadence flip: was 17:45 ET; moved to 19:30 ET so daily reports
# land late afternoon / early evening rather than mid-afternoon. Pre-flight
# runs at 18:30 (rexfinhub-preflight.timer) — 1h before this fires — to
# build previews + post the GO/HOLD summary. Gate open: 19:00. Gate close: 20:00.
OnCalendar=Mon..Fri 19:30:00
Persistent=true

[Install]
WantedBy=timers.target
