# deploy/systemd/rexfinhub-single-filing-worker.service
[Unit]
Description=REX FinHub Single-Filing Enricher Worker (Tier 2 — processes Tier 1 atom alerts)
After=network-online.target rexfinhub-atom-watcher.service
Wants=network-online.target

[Service]
Type=simple
User=jarvis
WorkingDirectory=/home/jarvis/rexfinhub
EnvironmentFile=/home/jarvis/rexfinhub/config/.env
# Pin process timezone to ET (audit fix R7).
Environment=TZ=America/New_York
Environment=SEC_CACHE_DIR=/home/jarvis/rexfinhub/cache/sec
Environment=SINGLE_FILING_BATCH_SIZE=20
Environment=SINGLE_FILING_POLL_INTERVAL=30
ExecStart=/home/jarvis/venv/bin/python -m etp_tracker.single_filing_worker
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
